# .config/systemd/user/waybar.service
[Unit]
Description=Waybar

[Service]
ExecStart=/usr/bin/waybar
Restart=no

[Install]
WantedBy=default.target

